# utils/build/virtual_machine/weblogs/python/test-app-python/test-app.service
[Unit]
Description=Python Weblog App Service
After=syslog.target network.target

[Service]
SuccessExitStatus=143

User=datadog

Type=simple

Environment=DD_APM_INSTRUMENTATION_DEBUG=TRUE
Environment=PYTHONUNBUFFERED=1
Environment=DJANGO_SETTINGS_MODULE=django_app
WorkingDirectory=/home/datadog
ExecStart=/bin/bash -c '/home/datadog/.pyenv/shims/python -m django runserver 0.0.0.0:5985 >> /home/datadog/app-std.out 2>&1 </dev/null'
ExecStop=/bin/kill -15 $MAINPID

[Install]
WantedBy=multi-user.target