# utils/build/virtual_machine/weblogs/nodejs/test-app-nodejs/test-app.service
[Unit]
Description=Nodejs Weblog App Service
After=syslog.target network.target

[Service]
SuccessExitStatus=143

User=datadog

Type=simple

Environment=DD_APM_INSTRUMENTATION_DEBUG=TRUE
Environment=DD_TRACE_DEBUG=true
WorkingDirectory=/home/datadog
ExecStart=/bin/bash -c 'node index.js >> /var/log/weblog/app.log 2>&1 </dev/null'
ExecStop=/bin/kill -15 $MAINPID

[Install]
WantedBy=multi-user.target