# utils/build/virtual_machine/weblogs/java/test-app-java/test-app.service
[Unit]
Description=Java Weblog App Service
After=syslog.target network.target

[Service]
SuccessExitStatus=143

User=datadog

Type=simple

Environment=DD_APM_INSTRUMENTATION_DEBUG=TRUE
Environment=DD_CONFIG_SOURCES=LOCAL:/etc/datadog-agent/inject/debug_config.yaml
WorkingDirectory=/home/datadog
ExecStart=/bin/bash -c 'java -Dserver.port=5985 -jar k8s-lib-injection-app-0.0.1-SNAPSHOT.jar >> /home/datadog/app-std.out 2>&1 </dev/null'
ExecStop=/bin/kill -15 $MAINPID

[Install]
WantedBy=multi-user.target